# Files/Other/rc.local
#!/bin/sh -e
iptables -t nat -A POSTROUTING -o ens3 -s 10.8.0.0/24 -j MASQUERADE
iptables -t nat -A POSTROUTING -s 10.8.0.0/24 -o eth0 -j MASQUERADE
badvpn-udpgw --listen-addr 127.0.0.1:7300 > /dev/nul &
 exit 0
